# 100block.service
[Unit]
Description=100 Block - What did you do?

[Service]
Type=oneshot
ExecStart=/home/jodie/100block/100block-prompt.sh
